# bootc/gateway-files/quadlets/pangolin.container
[Unit]
Description=Pangolin tunneled reverse proxy server
Documentation=https://fossorial.io/
After=pangolin-network.service
Wants=pangolin-network.service

[Container]
Image=docker.io/fosrl/pangolin:latest
ContainerName=pangolin
Network=pangolin-network
Volume=/var/lib/pangolin/config:/app/config:Z
HealthCmd=curl -f http://localhost:3001/api/v1/
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
AutoUpdate=registry
Label=org.opencontainers.image.title=pangolin
Label=org.opencontainers.image.description=Pangolin tunneled reverse proxy server

[Service]
Restart=unless-stopped
TimeoutStartSec=300
ExecStartPre=-/usr/bin/podman pull docker.io/fosrl/pangolin:latest

[Install]
WantedBy=multi-user.target